# com.henry.trading-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.henry.trading-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/henryschertzinger/Desktop/autonomous-agent/node_modules/.bin/tsx</string>
        <string>/Users/henryschertzinger/Desktop/autonomous-agent/agent-v3.1.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/henryschertzinger/Desktop/autonomous-agent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/henryschertzinger</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/henryschertzinger/Desktop/autonomous-agent/logs/agent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/henryschertzinger/Desktop/autonomous-agent/logs/agent-error.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
